# deploy_tools/gunicorn-upstart.template.conf
start on net-device-up
stop on shutdown

respawn

setuid ubuntu
chdir /home/ubuntu/www/SITENAME

exec ../../teamsurmandl/bin/gunicorn --bind unix:/tmp/SITENAME/.socket teamsurmandl.wsgi:application